# apps/r2-bridge/systemd/r2-camera-publisher.service
[Unit]
Description=RVEP R2 Camera Publisher (C270 → x264enc H.264 → LiveKit)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=mic-742
# 配合 Jetson 環境調整 Runtime
Environment=XDG_RUNTIME_DIR=/tmp/r2-camera-runtime
Environment=HOME=/tmp

# 修正 1：精準指向我們融入後的共用環境變數檔，且不加減號，確保強制讀取檢查
EnvironmentFile=/etc/rvep/r2-bridge.env

# 修正 2：修正為 Jetson 實際的 /usr/bin/mkdir 路徑
ExecStartPre=/usr/bin/mkdir -p /tmp/r2-camera-runtime

# 保持原有的 bash 封裝以確保環境變數正確傳遞
ExecStart=/bin/bash -c 'exec /opt/rvep/r2-bridge/r2_bridge/camera-publisher.sh'
Restart=on-failure
RestartSec=8
KillMode=mixed
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target