# apps/r2-bridge/systemd/r2-twist-mux.service
[Unit]
Description=RVEP twist_mux (priority merge → /cmd_vel)
After=network-online.target r2-wheeltec-base.service
Wants=network-online.target r2-wheeltec-base.service

[Service]
Type=simple
User=playrobot
Group=playrobot
Environment=ROS_DOMAIN_ID=0
ExecStart=/bin/bash -lc 'source /opt/ros/jazzy/setup.bash && exec ros2 run twist_mux twist_mux --ros-args --params-file /opt/rvep/r2-bridge/config/twist_mux.yaml -r cmd_vel_out:=/cmd_vel'
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
