# apps/r2-bridge/systemd/r2-wheeltec-base.service
[Unit]
Description=Wheeltec R2 base controller + LiDAR + TF (ROS2 jazzy launch)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=mic-742
Group=mic-742
Environment=ROS_DOMAIN_ID=0
Environment=HOME=/home/mic-742
# Wheeltec workspace install is on user's Desktop (colcon built-in paths).
WorkingDirectory=/home/mic-742
# Use a wrapper script in /opt/rvep — /tmp would not survive reboots.
ExecStart=/opt/rvep/r2-bridge/r2_bridge/wheeltec-base.sh
Restart=on-failure
# Give udev time to resolve /dev/huanyu_* symlinks after USB enumeration.
RestartSec=10
# launch.py spawns many child nodes — TERM cascades via process group.
KillMode=mixed
TimeoutStopSec=15
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
